# target/linux/econet/base-files/etc/uci-defaults/99-wps-led
#!/bin/sh
# Enable the WPS LED daemon (drives green:wps from hostapd WPS state).
/etc/init.d/wps-led enable
exit 0
